# s6/cron/run
#!/bin/sh

exec s6-setuidgid root /usr/sbin/crond -fS
